# system/udev/rules.d/99-battery-thresholds.rules
ACTION="change", SUBSYSTEM=="drm", RUN+="/usr/bin/systemctl start thunderbolt-charge-switch.service"
